# sambanas/rootfs/etc/s6-overlay/s6-rc.d/avahi/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start avahi service
# ==============================================================================
declare SMB_HOST

SMB_HOST=$(grep -i '^\s*netbios name\s*=' /etc/samba/smb.conf | cut -f2 -d= | tr -d '[:blank:]')

bashio::log.info "Starting the AVAHI for ${SMB_HOST%.*}..."

exec avahi-publish-service -v -f -s "${SMB_HOST%.*}" _smb._tcp 445
