# sambanas2/rootfs/etc/s6-overlay/s6-rc.d/modprobe/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# vim: ft=bash
# ==============================================================================
# Exec Modprobe for all fs
# ==============================================================================
bashio::log.debug  "-> Loading kernel modules..."

# Load necessary modules
MOD_FS=$(ls /lib/modules/$(uname -r)/kernel/fs)

for mod in $MOD_FS; do
    modprobe "$mod" ||:
done

fstypes=$(grep -v nodev </proc/filesystems | tr -d '\n')
bashio::log.blue "------------------------------------------------------------------------------------------------------"
bashio::log.green "System Supported fs: ${fstypes}"
if grep -q fuseblk </proc/filesystems; then bashio::log.green "Supported fusefs: $(find /usr/sbin -name "mount*" | cut -c 17- | tr "\n" " " | sed s/fuse.//g)"; fi
bashio::log.blue "------------------------------------------------------------------------------------------------------"

bashio::log.debug  "-> Done."
