# sambanas2/rootfs/etc/s6-overlay/s6-rc.d/modprobe/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# vim: ft=bash
# ==============================================================================
# Exec Modprobe for all fs
# ==============================================================================
bashio::log.debug  "-> Loading kernel modules..."

# Load necessary modules
MOD_FS=$(ls "/lib/modules/$(uname -r)/kernel/fs")

for mod in $MOD_FS; do
    modprobe "$mod" ||:
done

# ==============================================================================
# Download and load extra kernel modules from GitHub (experimental)
# ==============================================================================
if bashio::config.true 'use_external_kernel_modules'; then
    bashio::log.info "External kernel modules support enabled (experimental)..."

    HAOS_VERSION=$(bashio::os.version)
    BOARD=$(bashio::os.board)
    EXTRA_MODULES_DIR="/config/extra_modules"
    mkdir -p "${EXTRA_MODULES_DIR}"

    bashio::log.info "Looking for extra modules: haos_version=${HAOS_VERSION} board=${BOARD}"

    RELEASE_URL="https://api.github.com/repos/dianlight/hasos_more_modules/releases/tags/${HAOS_VERSION}"
    CURL_ERR=$(mktemp)
    if ! RELEASE_DATA=$(curl -sf --max-time 30 "${RELEASE_URL}" 2>"${CURL_ERR}"); then
        bashio::log.warning "No release found for HAOS ${HAOS_VERSION} in hasos_more_modules ($(cat "${CURL_ERR}" | head -1))"
        RELEASE_DATA=""
    fi
    rm -f "${CURL_ERR}"

    if [ -n "${RELEASE_DATA}" ]; then
        # Collect all asset URLs matching this HAOS version and board
        # Asset naming: {module_name}_{haos_version}_{board}.ko.xz
        MATCH_SUFFIX="_${HAOS_VERSION}_${BOARD}.ko.xz"
        MATCHED_COUNT=0

        while IFS= read -r asset_url; do
            asset_name=$(basename "${asset_url}")
            # Match assets ending with _{HAOS_VERSION}_{BOARD}.ko.xz
            if [[ "${asset_name}" == *"${MATCH_SUFFIX}" ]]; then
                # Derive the target .ko filename by stripping .xz
                ko_name="${asset_name%.xz}"
                ko_path="${EXTRA_MODULES_DIR}/${ko_name}"

                MATCHED_COUNT=$((MATCHED_COUNT + 1))

                if [ -f "${ko_path}" ]; then
                    bashio::log.info "Extra module already present, skipping download: ${ko_name}"
                else
                    bashio::log.info "Downloading extra kernel module: ${asset_name}"
                    xz_path="${EXTRA_MODULES_DIR}/${asset_name}"
                    CURL_DL_ERR=$(mktemp)
                    if curl -sfL --max-time 60 "${asset_url}" -o "${xz_path}" 2>"${CURL_DL_ERR}"; then
                        XZ_ERR=$(mktemp)
                        if xz -d "${xz_path}" 2>"${XZ_ERR}"; then
                            bashio::log.info "Decompressed: ${ko_name}"
                        else
                            bashio::log.warning "Failed to decompress ${asset_name}: $(cat "${XZ_ERR}" | head -1) - skipping"
                            rm -f "${xz_path}"
                        fi
                        rm -f "${XZ_ERR}"
                    else
                        bashio::log.warning "Failed to download ${asset_name} ($(cat "${CURL_DL_ERR}" | head -1)) - skipping"
                        rm -f "${xz_path}"
                    fi
                    rm -f "${CURL_DL_ERR}"
                fi

                # Load the module if the .ko is present
                if [ -f "${ko_path}" ]; then
                    mod_name="${ko_name%.ko}"
                    bashio::log.info "Loading extra kernel module: ${mod_name}"
                    MOD_ERR=$(mktemp)
                    insmod "${ko_path}" 2>"${MOD_ERR}" \
                        || modprobe "${mod_name}" 2>"${MOD_ERR}" \
                        || bashio::log.warning "Could not load module ${mod_name}: $(cat "${MOD_ERR}" | head -1) - skipping"
                    rm -f "${MOD_ERR}"
                fi
            fi
        done < <(echo "${RELEASE_DATA}" | jq -r '.assets[].browser_download_url' 2>/dev/null)

        if [ "${MATCHED_COUNT}" -eq 0 ]; then
            bashio::log.warning "No extra kernel modules found for haos_version=${HAOS_VERSION} board=${BOARD}"
        fi
    fi
fi

fstypes=$(grep -v nodev </proc/filesystems | tr -d '\n')
bashio::log.blue "------------------------------------------------------------------------------------------------------"
bashio::log.green "System Supported fs: ${fstypes}"
if grep -q fuseblk </proc/filesystems; then bashio::log.green "Supported fusefs: $(find /usr/sbin -name "mount*" | cut -c 17- | tr "\n" " " | sed s/fuse.//g)"; fi
bashio::log.blue "------------------------------------------------------------------------------------------------------"

bashio::log.debug  "-> Done."
