# sambanas2/rootfs/etc/s6-overlay/s6-rc.d/factory-reset/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# vim: ft=bash
# ==============================================================================
# Factory Reset - Delete SQL database and all configuration
# ==============================================================================
set -euo pipefail
exec 2>&1
trap 'bashio::log.error "factory-reset/run failed at line ${BASH_LINENO[0]} (exit code: ${?})"' ERR

# Check if factory_reset is enabled
if bashio::config.true 'factory_reset'; then
  bashio::log.blue "╔══════════════════════════════════════════════════════════════════════════════╗"
  bashio::log.red  "║                                                                              ║"
  bashio::log.red  "║                       ⚠️  FACTORY RESET IN PROGRESS  ⚠️                       ║"
  bashio::log.red  "║                                                                              ║"
  bashio::log.red  "║  This will DELETE all Samba NAS2 configurations, settings, and databases!    ║"
  bashio::log.red  "║                                                                              ║"
  bashio::log.blue "╚══════════════════════════════════════════════════════════════════════════════╝"
  
  bashio::log.warning "Starting factory reset process..."
  
  # Delete any database or related files (can be in /config or /data)
  bashio::log.warning "Cleaning up database files..."
  for db_dir in /config /data; do
    if [ -d "${db_dir}" ]; then
      rm -f "${db_dir}"/config.db3-* 2>/dev/null || true
      rm -f "${db_dir}"/*.db3 2>/dev/null || true
    fi
  done
  bashio::log.info "Database files cleaned"

  # Delete any certificate files
  bashio::log.warning "Cleaning up certificate files..."
  rm -f /config/server.{cert,key} 2>/dev/null || true
  bashio::log.info "Certificate files cleaned"


  # Delete downloaded extra kernel modules
  if [ -d "/data/extra_modules" ]; then
    bashio::log.warning "Removing downloaded extra kernel modules directory..."
    rm -rf /data/extra_modules
    bashio::log.info "Extra kernel modules directory removed"
  fi

  # Delete upgrade directory to remove any downloaded updates
  # Note: config/upgrade is used for development/testing, data/upgrade is used in production for downloaded updates
  if [ -d "/data/upgrade" ]; then
    bashio::log.warning "Removing upgrade directory to clean up any downloaded updates..."
    rm -rf /data/upgrade
    bashio::log.info "Upgrade directory removed"
  fi
  
  # Set clean_upgrade_dir to true so check-srat-update will clean it
  bashio::log.info "Setting clean_upgrade_dir to true for upgrade directory cleanup..."
  bashio::addon.option 'clean_upgrade_dir' true || \
    bashio::log.warning "Could not set clean_upgrade_dir to true"
  
  # Reset factory_reset flag to false
  bashio::log.info "Resetting factory_reset flag to false..."
  bashio::addon.option 'factory_reset' false || \
    bashio::log.warning "Could not reset factory_reset to false"
  
  bashio::log.blue  "╔══════════════════════════════════════════════════════════════════════════════╗"
  bashio::log.green "║                                                                              ║"
  bashio::log.green "║                      ✅  FACTORY RESET COMPLETED  ✅                         ║"
  bashio::log.green "║                                                                              ║"
  bashio::log.green "║  All configurations, settings, and databases have been deleted.              ║"
  bashio::log.green "║  The add-on will now start with default settings.                            ║"
  bashio::log.green "║                                                                              ║"
  bashio::log.blue  "╚══════════════════════════════════════════════════════════════════════════════╝"
else
  bashio::log.debug "Factory reset not requested, skipping..."
fi
