# sambanas2/rootfs/etc/s6-overlay/s6-rc.d/nmbd/run
#!/usr/bin/env bash
# ==============================================================================
# Start nmbd service
# ==============================================================================
# Add QUIC library to library path if it exists
if [ -d "/usr/local/quic/lib" ]; then
  export LD_LIBRARY_PATH="/usr/local/quic/lib:${LD_LIBRARY_PATH}"
fi

export PATH="/usr/local/samba/bin:/usr/local/samba/sbin:$PATH"


exec nmbd \
  --foreground \
  --debug-stdout \
  --no-process-group