# sambanas/rootfs/etc/services.d/wsdd-publish-service/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start avahi service
# ==============================================================================
declare SMB_GROUP
declare SMB_HOST


SMB_GROUP=$(grep -i '^\s*workgroup\s*=' /etc/samba/smb.conf | cut -f2 -d= | tr -d '[:blank:]')
SMB_HOST=$(grep -i '^\s*netbios name\s*=' /etc/samba/smb.conf | cut -f2 -d= | tr -d '[:blank:]')

bashio::log.info "Starting the WSDD daemon for ${SMB_GROUP}/${SMB_HOST}..."

# Get supported interfaces
for interface in $(bashio::network.interfaces); do
    interfaces+=("-i ${interface}")
done

if [ ${#interfaces[@]} -eq 0 ]; then
    bashio::exit.nok 'No supported interfaces found to bind on.'
fi
bashio::log.info "Interfaces: $(printf '%s ' "${interfaces[@]}")"

exec wsdd -v $(printf '%s ' "${interfaces[@]}") -n ${SMB_HOST} -w ${SMB_GROUP}