# sambanas/rootfs/etc/services.d/mqtt-handler-service/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start mqtt service
# ==============================================================================
if [ -f /root/.config/mosquitto_pub ]; then
    bashio::log.info "Starting the MQTT daemon for disks info..."
    exec /usr/local/bin/mqtt_daemon.sh
else
    bashio::log.warning "No MQTT config."
    exec sleep 1440
fi    