# sambanas2/rootfs/etc/s6-overlay/s6-rc.d/nfsd/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# vim: ft=bash
# ============================================================================
# Start NFS server stack under s6 supervision
# Exports are automatically managed by SRAT for HA storage directories
# ============================================================================
set -euo pipefail

EXPORTS_FILE="${EXPORTS_FILE:-/etc/exports}"
NFS_THREADS="${NFS_THREADS:-8}"
NFS_ENABLE_VERSIONS="${NFS_ENABLE_VERSIONS:-4 4.1 4.2}"
NFS_DISABLE_VERSIONS="${NFS_DISABLE_VERSIONS:-3}"
STATD_PID_FILE="/run/nfsd.statd.pid"
IDMAPD_PID_FILE="/run/nfsd.idmapd.pid"
NFSDCLD_PID_FILE="/run/nfsd.nfsdcld.pid"

build_nfsd_version_flags() {
  NFS_VERSION_FLAGS=()

  for version in ${NFS_DISABLE_VERSIONS}; do
    NFS_VERSION_FLAGS+=("-N" "${version}")
  done

  for version in ${NFS_ENABLE_VERSIONS}; do
    NFS_VERSION_FLAGS+=("-V" "${version}")
  done
}

mountd_debug_level() {
  local level

  level="$(bashio::config 'log_level' 'info')"
  case "${level}" in
    trace) echo 7 ;;
    debug) echo 6 ;;
    info) echo 5 ;;
    notice) echo 4 ;;
    warning) echo 3 ;;
    error) echo 2 ;;
    fatal) echo 1 ;;
    *)
      bashio::log.warning "rpc.mountd: Unknown log_level '${level}'; defaulting to warning."
      echo 3
      ;;
  esac
}

ensure_rpcbind() {
  if pgrep -x rpcbind >/dev/null; then
    return
  fi

  bashio::log.info "Starting rpcbind for NFS services"
  rpcbind -w
}

start_statd() {
  if pgrep -x rpc.statd >/dev/null; then
    return
  fi

  bashio::log.info "Starting rpc.statd"
  rpc.statd -F --no-notify &
  echo "$!" >"${STATD_PID_FILE}"
}

start_idmapd() {
  if ! command -v rpc.idmapd >/dev/null; then
    return
  fi

  if pgrep -x rpc.idmapd >/dev/null; then
    return
  fi

  bashio::log.info "Starting rpc.idmapd"
  rpc.idmapd -f -p /run/rpc_pipefs  &
  echo "$!" >"${IDMAPD_PID_FILE}"
}

start_nfsdcld() {
  if ! command -v nfsdcld >/dev/null; then
    return
  fi

  if pgrep -x nfsdcld >/dev/null; then
    return
  fi

  bashio::log.info "Starting nfsdcld"
  nfsdcld -F  -p /run/rpc_pipefs &
  echo "$!" >"${NFSDCLD_PID_FILE}"
}

if [ ! -s "${EXPORTS_FILE}" ]; then
  bashio::log.info "NFS exports file ${EXPORTS_FILE} not yet generated by SRAT; waiting for configuration."
fi

while [ ! -s "${EXPORTS_FILE}" ]; do
  sleep 30
done

ensure_rpcbind
start_statd
start_idmapd
start_nfsdcld

build_nfsd_version_flags

while true; do
  bashio::log.info "Starting rpc.nfsd (threads=${NFS_THREADS}) enable=[${NFS_ENABLE_VERSIONS}] disable=[${NFS_DISABLE_VERSIONS}]"
  bashio::log.debug "NFS version flags: ${NFS_VERSION_FLAGS[@]}"

  if ! rpc.nfsd "${NFS_THREADS}" "${NFS_VERSION_FLAGS[@]}"; then
    bashio::log.error "rpc.nfsd failed to start (threads=${NFS_THREADS}); retrying in 10s."
    sleep 10
    continue
  fi

  if ! exportfs -r; then
    bashio::log.error "Failed to apply exports from ${EXPORTS_FILE}; retrying in 10s."
    rpc.nfsd 0 || true
    sleep 10
    continue
  fi

  bashio::log.info "NFS services started using ${EXPORTS_FILE}"
  MOUNTD_DEBUG_LEVEL="$(mountd_debug_level)"
  bashio::log.debug "rpc.mountd debug level: ${MOUNTD_DEBUG_LEVEL}"
  exec s6-notifyoncheck -c "showmount -e localhost >/dev/null 2>&1" \
    rpc.mountd -F --manage-gids -d "${MOUNTD_DEBUG_LEVEL}"
done