# sambanas2/rootfs/etc/s6-overlay/s6-rc.d/nfsd/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# vim: ft=bash
# ============================================================================
# Start NFS server stack under s6 supervision
# Exports are automatically managed by SRAT for HA storage directories
# ============================================================================
set -uo pipefail

EXPORTS_FILE="${EXPORTS_FILE:-/etc/exports}"
NFS_THREADS="${NFS_THREADS:-8}"
NFS_ENABLE_VERSIONS="${NFS_ENABLE_VERSIONS:-4 4.1 4.2}"
NFS_DISABLE_VERSIONS="${NFS_DISABLE_VERSIONS:-3}"

build_nfsd_version_flags() {
  NFS_VERSION_FLAGS=()
  for version in ${NFS_DISABLE_VERSIONS}; do
    NFS_VERSION_FLAGS+=("-N" "${version}")
  done
  for version in ${NFS_ENABLE_VERSIONS}; do
    NFS_VERSION_FLAGS+=("-V" "${version}")
  done
}

# Ensure rpcbind is running (required by rpc.mountd and rpc.statd)
ensure_rpcbind() {
  if pgrep -x rpcbind >/dev/null 2>&1; then
    return 0
  fi
  bashio::log.info "Starting rpcbind"
  # -w: warm start (reads existing state); -f: foreground fallback
  if ! rpcbind 2>/dev/null; then
    bashio::log.warning "rpcbind daemon mode failed; trying foreground mode"
    rpcbind -f &
  fi
  sleep 1
}

# rpc.statd: Network Status Monitor — needed for NFS lock recovery
ensure_statd() {
  if pgrep -x rpc.statd >/dev/null 2>&1; then
    return 0
  fi
  # Ensure state directory exists (statd will fail without it)
  mkdir -p /var/lib/nfs/statd/sm /var/lib/nfs/statd/sm.bak
  bashio::log.info "Starting rpc.statd"
  # -F: foreground  -L: suppress reboot notifications (no-notify)
  rpc.statd -F -L &
}

# rpc.idmapd: NFSv4 ID mapping daemon (required for NFSv4)
ensure_idmapd() {
  if ! command -v rpc.idmapd >/dev/null 2>&1; then
    return 0
  fi
  if pgrep -x rpc.idmapd >/dev/null 2>&1; then
    return 0
  fi
  bashio::log.info "Starting rpc.idmapd"
  # -f: foreground  -p: rpc_pipefs directory (mounted by init-nfsd)
  rpc.idmapd -f -p /run/rpc_pipefs &
}

# nfsdcld: NFSv4.0 client tracking daemon (recommended for crash recovery)
ensure_nfsdcld() {
  if ! command -v nfsdcld >/dev/null 2>&1; then
    return 0
  fi
  if pgrep -x nfsdcld >/dev/null 2>&1; then
    return 0
  fi
  bashio::log.info "Starting nfsdcld"
  # -F: foreground  -p: rpc_pipefs directory  -s: client state storage
  nfsdcld -F -p /run/rpc_pipefs -s /var/lib/nfs/v4recovery &
}

# Start or restart all NFS support daemons
ensure_support_daemons() {
  ensure_rpcbind
  ensure_statd
  ensure_idmapd
  ensure_nfsdcld
}

# Wait for SRAT to generate the exports file
if [ ! -s "${EXPORTS_FILE}" ]; then
  bashio::log.info "Waiting for NFS exports file (${EXPORTS_FILE}) to be generated by SRAT..."
fi
while [ ! -s "${EXPORTS_FILE}" ]; do
  sleep 30
done

build_nfsd_version_flags

# Initial daemon startup
ensure_support_daemons

while true; do
  # Re-check and restart any crashed support daemons before each nfsd launch
  ensure_support_daemons

  bashio::log.info "Starting rpc.nfsd (threads=${NFS_THREADS}) enable=[${NFS_ENABLE_VERSIONS}] disable=[${NFS_DISABLE_VERSIONS}]"

  if ! rpc.nfsd "${NFS_THREADS}" "${NFS_VERSION_FLAGS[@]}"; then
    bashio::log.error "rpc.nfsd failed to start; retrying in 10s."
    sleep 10
    continue
  fi

  if ! exportfs -r; then
    bashio::log.error "Failed to apply exports from ${EXPORTS_FILE}; retrying in 10s."
    rpc.nfsd 0 || true
    sleep 10
    continue
  fi

  bashio::log.info "NFS server started — exports applied from ${EXPORTS_FILE}"

  # Execute rpc.mountd as the supervised foreground process.
  # Health check: verify that the kernel nfsd has at least one active thread.
  # (Does not depend on rpcbind/portmapper, so it works for NFSv4-only setups.)
  exec s6-notifyoncheck -c "grep -qv '^0$' /proc/fs/nfsd/threads 2>/dev/null" \
    rpc.mountd -F --manage-gids
done
