# sambanas/rootfs/etc/s6-overlay/s6-rc.d/init-mqtt/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# MQTT autodiscovery cleanup
# ==============================================================================
declare topic

if ! bashio::config.true "autodiscovery.disable_autoremove"; then
    bashio::log.info "MQTT cleanup."
    topic=$(bashio::config 'mqtt_topic')
    if [ "$topic" = "null" ]; then topic="sambanas"; fi
    mosquitto_sub -t "homeassistant/sensor/${topic}/+/config" --remove-retained -W 5 -R -E
    bashio::log.info "MQTT cleanup Done."
fi
