# sambanas2/rootfs/etc/s6-overlay/s6-rc.d/check-srat-update/up
/etc/s6-overlay/s6-rc.d/check-srat-update/run