# besim/rootfs/etc/s6-overlay/s6-rc.d/besim/run
#!/usr/bin/env bash
# ==============================================================================
# Start besim service
# ==============================================================================
exec /usr/bin/python3 /opt/BeSIM/app.py