# besim/rootfs/etc/s6-overlay/s6-rc.d/besim/run
#!/command/with-contenv bashio
# ==============================================================================
# Start besim service
# ==============================================================================
log_level=$(bashio::string.lower "$(bashio::config log_level info)")
args="-l ${log_level^^}"
if bashio::config.true 'work_as_proxy'; then
    args="${args} -p $(bashio::config upstream_dns)"
fi
bashio::log.info ${args}
exec python3 /opt/BeSIM/app.py ${args}
