# sambanas/rootfs/etc/s6-overlay/s6-rc.d/hd-idle/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start hd-idle service
# ==============================================================================
if bashio::config.has_value 'hdd_idle_seconds'; then
    bashio::log.info "Enabling HDD IDLE after $(bashio::config 'hdd_idle_seconds')sec"
    exec hd-idle -i "$(bashio::config 'hdd_idle_seconds')"
else
    exec sleep infinity
fi
