# sambanas/rootfs/etc/s6-overlay/s6-rc.d/srat/run
#!/usr/bin/env bash
# ==============================================================================
# Start nmbd service
# ==============================================================================
exec /usr/local/bin/srat_$(arch) \
  -port 3000 \
  -conf /config/bootconfig.json
