# besim/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bash

/usr/sbin/nginx -g "daemon off;"
