# sambanas2/rootfs/etc/s6-overlay/s6-rc.d/clean-config/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# vim: ft=bash
# ==============================================================================
# Remove deprecated config options 
# ==============================================================================
bashio::log.debug  "-> Clean config ..."
if bashio::config.has_value 'enable_smart'; then
    bashio::log.warning "Option 'enable_smart' is deprecated and will be removed from the config."
    bashio::addon.option 'enable_smart'
fi
if bashio::config.has_value 'hdd_idle_seconds'; then
    bashio::log.warning "Option 'hdd_idle_seconds' is deprecated and will be removed from the config."
    bashio::addon.option 'hdd_idle_seconds'
fi
bashio::log.debug  "-> Done."

