# sambanas2/rootfs/etc/s6-overlay/s6-rc.d/smbd/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# vim: ft=bash
# ==============================================================================
# Start smbd service
# ==============================================================================

# Add QUIC library to library path if it exists
if [ -d "/usr/local/quic/lib" ]; then
  export LD_LIBRARY_PATH="/usr/local/quic/lib:${LD_LIBRARY_PATH}"
fi

export PATH="/usr/local/samba/bin:/usr/local/samba/sbin:$PATH"

#exec smbd \
#  --foreground \
#  --debug-stdout \
#  --no-process-group

# exec ensures s6-notifyoncheck replaces the shell process
# -c: The command to check health (smbclient listing shares)
# smbd -F: Runs in foreground (mandatory for s6)
exec s6-notifyoncheck -c "smbclient -L localhost -N" \
    smbd -F --no-process-group --debug-stdout <&-  